# systemd-units/emoney.service
[Unit]
Description=E-Money Full Node

[Service]
Environment=EMD_P2P_LADDR=tcp://0.0.0.0:2140
Environment=EMD_RPC_LADDR=tcp://0.0.0.0:2141
Environment=EMD_NODE=tcp://127.0.0.1:2141
Environment=EMD_GRPC_ADDRESS=0.0.0.0:2142
Environment=EMD_API_ADDRESS=tcp://127.0.0.1:2143
Environment=EMD_GRPC_WEB_ADDRESS=127.0.0.1:2144
Environment=EMD_P2P_PERSISTENT_PEERS="6fe3b41ddaf1e5d7191bd09621c6c0412a118c64@139.59.8.62:2140,0ad7bc7687112e212bac404670aa24cd6116d097@50.18.83.75:26656,1723e34f45f54584f44d193ce9fd9c65271ca0b3@13.124.62.83:26656,ecec8933d80da5fccda6bdd72befe7e064279fc1@207.180.213.123:26676,34eca4a9142bf9c087a987b572c114dad67a8cc5@172.105.148.191:26656,0b186517e4d82eb4c000a567e486b7b96bf19752@44.195.95.22:26656,eed66085c975189e3d498fe61af2fcfb3da34924@217.79.184.40:26656,226d713acb10a1cd09b3b890406b886d7b071cf2@10.105.2.101:26656,552cb87c40f879c8c41005eeb36030502671e7a4@18.156.91.213:26656,170255a4a9b120786d2d6c3931ebf0125087cca0@62.171.171.139:26656,d12ac664ff9766e0535c836173111e3eed462e79@185.144.83.157:26656,d920b1fc64976f5838bc2bc7924105c177f1e526@88.198.62.198:2140"
Environment=EMD_P2P_UNCONDITIONAL_PEERS="c2766f7f6dfe95f2eb33e99a538acf3d6ec608b1"
Environment=EMD_P2P_SEEDS="17533be7691494a0bc2dbd174930fa6a57b5b98c@3.69.44.220:30656,3f6ce97077a9bddd6dd637cd64c140eabdebbda5@3.68.159.237:30656,170255a4a9b120786d2d6c3931ebf0125087cca0@62.171.171.139:26656"
Environment=EMD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=EMD_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=EMD_PRUNING=nothing
User=root
Restart=always
RestartSec=30
ExecStart=/root/go/bin/emd start

[Install]
WantedBy=multi-user.target
