# systemd-units/osmosis.service
[Unit]
Description=Osmosis Full Node

[Service]
Environment=OSMOSISD_P2P_LADDR=tcp://0.0.0.0:2000
Environment=OSMOSISD_RPC_LADDR=tcp://0.0.0.0:2001
Environment=OSMOSISD_NODE=tcp://127.0.0.1:2001
Environment=OSMOSISD_GRPC_ADDRESS=0.0.0.0:2002
Environment=OSMOSISD_API_ADDRESS=tcp://127.0.0.1:2003
Environment=OSMOSISD_GRPC_WEB_ADDRESS=127.0.0.1:2004
Environment=OSMOSISD_P2P_UNCONDITIONAL_PEERS="63aba59a7da5197c0fbcdc13e760d7561791dca8"
Environment=OSMOSISD_P2P_SEED_NODES=1b077d96ceeba7ef503fb048f343a538b2dcdf1b@136.243.218.244:26656,2308bed9e096a8b96d2aa343acc1147813c59ed2@3.225.38.25:26656,902bdfe51b6a97cc9369664a21c87ed61d471d2a@136.243.218.243,085f62d67bbf9c501e8ac84d4533440a1eef6c45@95.217.196.54:26656,f515a8599b40f0e84dfad935ba414674ab11a668@osmosis.blockpane.com:26656
Environment=OSMOSISD_P2P_PERSISTENT_PEERS="204a5b8f1e263ee8dc64942aba0fc22e39c182b5@135.125.104.217:34656,4124dec2db38c7fd7564c819d7921e8a679af9ce@47.240.17.20:29656,6a2c5c53ee600a8876f2a2df5c5858018cb04a98@143.198.190.100:26656,6caae003ee30049a2f39cebf2cfdc0bf71eb118a@5.161.44.25:26656,b1bc0463a4a5319266703df922735996c14d3cbf@104.149.131.242:26656,5e3001ed7a0362dca98849d460bb81406d30ce5e@147.75.199.219:30441,a3ce433431b8a82bdc06ee7e24b62f062c189fe3@88.99.105.146:26656,8dcbb97b7f61fed4d6605aae3123f2e6764dca23@95.179.190.79:26656,f58675dc341ddc0a32fbf90771d77c7467995f6c@18.193.72.160:29656,41c51383aadcd4bb1701f6192cc6fd0fa7caa135@147.75.195.129:26656,56d5b2b168a5d6a727c05ea1591353c5aa183958@139.162.38.11:26656,d264c4d76a4634e622b3f85754110ae76899b155@15.222.251.55:26711,acf30e4e2e7d72e8dfa36f26d479a854143f25ab@51.210.154.28:34656,f82d7f54601f46cc7aded1d27deff872e85733da@23.88.81.42:2000,48aa143ccbef5c487fbd26602c64e7caa26d1668@45.63.120.82:26656,457164de65c1c01f32ec2a4a4738709362ac9a3d@134.122.35.19:26656,87b1f4b92480ad69f1789041665c2ebd5e894e01@194.233.71.195:26656,8f67a2fcdd7ade970b1983bf1697111d35dfdd6f@52.79.199.137:26656,e437756a853061cc6f1639c2ac997d9f7e84be67@144.76.183.180:26656,febc4e06bb4a42782c4a97e6992515237a00216e@135.125.8.184:26656,c05b2db11824b31d97c9b12d857ac45ee537afb0@137.184.31.73:26656,19004ec3af2a6d439ddfc61a90adc381b074b160@65.21.249.233:26656,f2ad9f2b84864ad11e122a90911268893fd724c6@3.67.75.111:29656,70866c14b54ab13632e433d291fbbb0ef228e467@162.0.227.134:26656,b7f2ffa619cbfce0d6a40c1a09a74cf45d86f4bb@52.9.118.174:26656,710de9ec42627fe25ce02613d545e7aa92c04ea5@147.75.195.129:26656,d943fea86e04b98303e65b273ec026f9c5924b26@49.12.233.51:26656,079f13034b942f8b7e9fa98d8c4c2b3bfddd71d8@18.156.4.168:26656,499b17bd696172a45f53f9b27763bcac83c83094@159.69.214.201:26656,11dfe200894f38e411beca77928e9dd118e66813@94.130.98.157:26656,f1ec480b42167e45ab5ee300940364b316b4d2f9@135.181.217.231:26656,f7236c3f33c562cd5afa9557bc3c0531ed8e99d6@35.183.177.221:26656,d3d917b91969e2914049730b6c79f0cc1f98d7fb@95.216.148.25:26656,4a23efbd954d7498e3018a21dfc9b7ab19d5f881@157.90.129.103:15609,ddefd4572c7ed6831b25d6f8b34833e2608df76d@35.86.22.120:26656,d9d451a3c051dc5d0d1d16d370897a29cf9535f1@44.234.147.28:26656,dd371acd778846afa33f5138c376da493d0f56cc@147.75.199.219:31970,83d6bf3da8c84072de73088ebdac9d0a1b7b08ad@49.12.156.37:26656,db78a2b79c12d57fa299dd1b1d9959ce2d6b2a8d@3.121.77.150:26656,6496bfe0d5016d2e6c9741d4894b0fb01ce7b050@52.79.110.125:26656,ace263b6a0ce646537b582f5b9fe8903cd05ed93@149.202.46.150:26656,53fb7c1e3739d27dc9710551c40717e071c03e05@5.9.230.132:26656,ffd72a2c9672d6c7d29d6f94afea95ac338ba6dc@176.34.75.138:26656,db6fcb90ba39d35f207e68588acb19ba2c7bb3e0@77.87.108.21:26656,f4e95eaa94302041aec1f24dcad6296c64efcb14@51.15.63.113:26656,11bfced94d58357917fe9189fb34d3def790b454@139.59.8.62:2000"
Environment=OSMOSISD_P2P_MAX_NUM_INBOUND_PEERS=1000
Environment=OSMOSISD_P2P_SEED_MODE=true
Environment=OSMOSISD_P2P_MAX_NUM_OUTBOUND_PEERS=200
Environment=OSMOSISD_PRUNING=nothing
Environment=OSMOSISD_HOME=/osmo1/
LimitNOFILE=500000
User=root
Restart=always
ExecStart=/root/go/bin/osmosisd start

[Install]
WantedBy=multi-user.target
